# images/logroom-ui/supervisord.conf
[supervisord]
nodaemon=true

[program:logroom-ui]
directory = /opt/logRoomUi
command=bash -c 'npm start > /var/log/logroom/ui 2>&1'
